# vmenuhelper/com.brianshumate.vmenu.helper.plist
<?xml version="1.0" encoding="UTF-8" ?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent plist for the vmenu XPC helper.

  This file is embedded in the main app bundle at:
    Contents/Library/LaunchAgents/com.brianshumate.vmenu.helper.plist

  SMAppService.agent(plistName:) reads this to register the helper with
  launchd. The MachServices key exposes the XPC endpoint that the
  sandboxed main app connects to.

  macOS 26 (Tahoe) notes:
  - KeepAlive with SuccessfulExit=false and Crashed=true ensures launchd
    restarts the helper after crashes or unexpected exits.
  - ProcessType=Interactive is recommended for XPC services that respond
    to user actions to get appropriate scheduling priority.
  - ThrottleInterval limits restart rate if the helper keeps crashing.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brianshumate.vmenu.helper</string>
    <key>BundleProgram</key>
    <string>Contents/MacOS/com.brianshumate.vmenu.helper</string>
    <key>MachServices</key>
    <dict>
        <key>com.brianshumate.vmenu.helper</key>
        <true />
    </dict>
    <key>RunAtLoad</key>
    <false />
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false />
        <key>Crashed</key>
        <true />
    </dict>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
